# legacy-ansible/stage1-jumpbox/roles/routing/templates/nested-lab-nft.service.j2
# Managed by nested-vsphere-lab (routing).
[Unit]
Description=Nested vSphere Lab nftables ruleset
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/nft -f {{ lab.state_dir }}/nftables.conf
ExecReload=/usr/sbin/nft -f {{ lab.state_dir }}/nftables.conf
ExecStop=/usr/sbin/nft flush ruleset

[Install]
WantedBy=multi-user.target
